# munin/trunk/munin-html.socket
[Unit]
Description=Munin FastCGI HTML Socket

[Socket]
ListenStream=/run/munin/fcgi-html.sock
Accept=yes

[Install]
WantedBy=sockets.target
