# repos/core-x86_64/gpm.path
[Unit]
Description=Virtual console mouse server

[Path]
PathExists=/dev/input/mice
