# btrfs-progs/repos/testing-i686/btrfs-scrub@.timer
[Unit]
Description=Monthly Btrfs scrub on %I

[Timer]
OnCalendar=monthly
AccuracySec=1d
Persistent=true

[Install]
WantedBy=multi-user.target
