# postfix/repos/testing-x86_64/postfix.service
[Unit]
Description=Postfix Mail Transport Agent
After=network.target

[Service]
CapabilityBoundingSet=~ CAP_NET_ADMIN CAP_SYS_ADMIN CAP_SYS_BOOT CAP_SYS_MODULE
ExecReload=/usr/bin/postfix reload
ExecStart=/usr/bin/postfix start
ExecStop=/usr/bin/postfix stop
PIDFile=/var/spool/postfix/pid/master.pid
PrivateDevices=true
PrivateTmp=true
ProtectSystem=true
Restart=always
Type=forking

[Install]
WantedBy=multi-user.target
