# inetutils/repos/core-i686/ftpd.service
[Unit]
Description=FTPD Daemon
After=network.target

[Service]
ExecStart=/usr/bin/ftpd -D
Type=forking

[Install]
WantedBy=multi-user.target
