# btrfs-progs/repos/core-i686/70-btrfs.rules
ACTION!="remove", SUBSYSTEM=="block", ENV{ID_FS_TYPE}=="btrfs", RUN+="/usr/bin/btrfs device scan /dev/%k"
