# btrfs-progs/trunk/btrfs-scrub@.service
[Unit]
Description=Btrfs scrub on %I

[Service]
Nice=19
IOSchedulingClass=idle
ExecStart=/usr/bin/btrfs scrub start -B %I
