# ulogd/repos/extra-x86_64/ulogd.service
[Unit]
Description=Netfilter Ulogd daemon
After=network.target

[Service]
ExecStart=/usr/bin/ulogd

[Install]
WantedBy=multi-user.target
