# man-db/repos/core-x86_64/man-db.timer
[Unit]
Description=Daily man-db cache purge

[Timer]
OnCalendar=daily
AccuracySec=1d
Persistent=true
