# lm_sensors/repos/extra-i686/fancontrol.service
[Unit]
Description=Fan control daemon

[Service]
PIDFile=/var/run/fancontrol.pid
ExecStart=/usr/bin/fancontrol

[Install]
WantedBy=multi-user.target
