# repos/core-i686/adsl.service
[Unit]
Description=ADSL Deamon

[Service]
Type=forking
ExecStart=/usr/bin/pppoe-start
ExecStop=/usr/bin/pppoe-stop

[Install]
WantedBy=multi-user.target
