# xsp/repos/extra-i686/xsp.service
[Unit]
Description=XSP Web Server

[Service]
ExecStart=/usr/bin/xsp --appconfigdir /etc/xsp --nonstop

[Install]
WantedBy=multi-user.target
