# man/trunk/whatis.cron.daily
#!/bin/sh

# Update the "whatis" database
/usr/sbin/makewhatis -u -w

