# repos/testing-any/pkgstats.cron
#!/bin/sh
sudo -n -u nobody -- /usr/bin/pkgstats -q
