# trunk/dhcpcd.service
[Unit]
Description=dhcpcd on all interfaces
Wants=network.target
Before=network.target

[Service]
Type=forking
PIDFile=/run/dhcpcd.pid
ExecStart=/bin/dhcpcd -q -b
ExecStop=/bin/dhcpcd -x

[Install]
WantedBy=multi-user.target
