# trunk/gpm.service
[Unit]
Description=Virtual console mouse server

[Service]
Type=forking
ExecStart=/usr/sbin/gpm

[Install]
WantedBy=multi-user.target
