# trunk/pkgstats.timer
[Unit]
Description=Weekly pkgstats submission

[Timer]
OnCalendar=weekly
AccuracySec=1d
Persistent=true
