# repos/testing-x86_64/bftpd.service
[Unit]
Description=bftpd daemon

[Service]
Type=forking
ExecStart=/usr/bin/bftpd -d
GuessMainPID=yes

[Install]
WantedBy=multi-user.target
