# repos/testing-x86_64/edac.service
[Unit]
Description=Initialize EDAC Drivers For Machine Hardware

[Service]
Type=oneshot
ExecStart=/usr/bin/edac-ctl --register-labels
ExecStop=/usr/bin/edac-ctl --unload
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
