# repos/testing-x86_64/krb5-kpropd.service
[Unit]
Description=Kerberos 5 propagation server

[Service]
ExecStart=/usr/bin/kpropd -S

[Install]
WantedBy=multi-user.target
