# repos/testing-x86_64/krb5-kadmind.service
[Unit]
Description=Kerberos 5 administration server

[Service]
ExecStart=/usr/bin/kadmind -nofork

[Install]
WantedBy=multi-user.target
