# repos/testing-x86_64/mdadm.service
[Unit]
Description=MDADM Event Monitor

[Service]
ExecStart=/sbin/mdadm --monitor --scan

[Install]
WantedBy=multi-user.target
