# repos/testing-x86_64/rlogin@.service
[Unit]
Description=Remote Login Facilities Server
After=local-fs.target

[Service]
ExecStart=-/usr/bin/rlogind
StandardInput=socket
