# repos/testing-x86_64/rsh.socket
[Unit]
Description=Remote Shell Facilities Activation Socket

[Socket]
ListenStream=514
Accept=true

[Install]
WantedBy=sockets.target
