# repos/testing-x86_64/telnet@.service
[Unit]
Description=Telnet Server
After=local-fs.target

[Service]
ExecStart=-/usr/bin/telnetd
StandardInput=socket
