# repos/testing-x86_64/updatedb.timer
[Unit]
Description=Daily locate database update

[Timer]
OnCalendar=daily
RandomizedDelaySec=12h
Persistent=true
