# repos/testing-x86_64/usbmuxd.service
[Unit]
Description=iOS USB Muxer

[Service]
ExecStart=/usr/sbin/usbmuxd -u -f -U usbmux
