# at/repos/extra-x86_64/atd.service
[Unit]
Description=ATD daemon

[Service]
Type=forking
ExecStart=/usr/sbin/atd

[Install]
WantedBy=multi-user.target
