# tftp-hpa/repos/extra-i686/tftpd.socket
[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target
