# tftp-hpa/repos/staging-x86_64/tftpd.socket
[Unit]
Description=hpa's original TFTP daemon
Conflicts=tftpd.service

[Socket]
ListenDatagram=69

[Install]
WantedBy=sockets.target
