# shadow/repos/core-x86_64/shadow.cron.daily
#!/bin/sh

# Verify integrity of password and group files
/usr/bin/pwck -r
/usr/bin/grpck -r

