# man/repos/core-i686/whatis.cron.daily
#!/bin/sh

# Make sure MANPATH is complete
source /etc/profile
source /etc/profile.d/* 2>/dev/null

# Update the "whatis" database
/usr/sbin/makewhatis -u -w

