# timidity++/repos/staging-x86_64/timidity.service
[Unit]
Description=TiMidity++ Daemon
After=sound.target

[Service]
ExecStart=/usr/bin/timidity -iA

[Install]
WantedBy=multi-user.target
