# radvd/repos/community-x86_64/radvd.service
[Unit]
Description=IPv6 Router Advertisement Daemon
After=network.target

[Service]
ExecStart=/usr/bin/radvd --nodaemon --logmethod=stderr

[Install]
WantedBy=multi-user.target
