# mlocate/trunk/updatedb.timer
[Unit]
Description=Daily locate database update

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d
