# hylafax/repos/staging-x86_64/hfaxd.service
[Unit]
Description=Hyla FAX client-server

[Service]
Type=forking
ExecStart=/usr/lib/fax/hfaxd -i hylafax
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
