# ipset/repos/extra-x86_64/ipset.service
[Unit]
Description=Loading IP Sets
Before=iptables.service ip6tables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/ipset -f /etc/ipset.conf restore
ExecReload=/usr/bin/ipset -f /etc/ipset.conf restore
ExecStop=/usr/bin/ipset destroy

[Install]
WantedBy=multi-user.target
