# bftpd/repos/extra-x86_64/bftpd_AT.service
[Unit]
Description=bftpd daemon

[Service]
ExecStart=/usr/bin/bftpd -i
StandardInput=socket
StandardOutput=socket
StandardError=socket
