# inetutils/repos/testing-x86_64/rlogin.socket
[Unit]
Description=Remote Login Facilities Activation Socket

[Socket]
ListenStream=513
Accept=true

[Install]
WantedBy=sockets.target
