# inetutils/repos/testing-x86_64/rsh@.service
[Unit]
Description=Remote Shell Facilities Server
After=local-fs.target

[Service]
ExecStart=-/usr/bin/rshd
StandardInput=socket
