# trunk/sources/iana-versioned.timer
[Unit]
Description=Periodically version IANA files

[Install]
WantedBy=default.target

[Timer]
OnCalendar=05:10
Persistent=true
