# trunk/sources/iana-versioned.service
[Unit]
Description=Version IANA files

[Install]
WantedBy=default.target

[Service]
ExecStart=/home/bisson/bin/iana-versioned.sh
