# inetutils/repos/core-x86_64/ftpd.service
[Unit]
Description=FTPD Daemon

[Service]
ExecStart=/usr/sbin/ftpd -D
Type=forking

[Install]
WantedBy=multi-user.target
