# openldap/repos/testing-i686/slapd.service
[Unit]
Description=OpenLDAP server daemon

[Service]
Type=forking
ExecStart=/usr/sbin/slapd -u ldap -g ldap

[Install]
WantedBy=multi-user.target
