# rp-pppoe/repos/core-x86_64/adsl.service
[Unit]
Description=ADSL Deamon

[Service]
Type=forking
ExecStart=/usr/sbin/pppoe-start
ExecStop=/usr/sbin/pppoe-stop

[Install]
WantedBy=multi-user.target
