# repos/extra-any/pkgstats.service
[Unit]
Description=Submit pkgstats

[Service]
Type=oneshot
ExecStart=/usr/bin/pkgstats -q
User=nobody
Group=nobody
