# repos/extra-any/pkgstats.cron
#!/bin/sh
su -s /bin/bash -c '/usr/bin/pkgstats -q' nobody
