# git/trunk/git-daemon.socket
[Unit]
Description=Git Daemon Socket
After=network.target

[Socket]
ListenStream=9418
Accept=true

[Install]
WantedBy=sockets.target
