# repos/testing-i686/mdadm.service
[Unit]
Description=MDADM Event Monitor

[Service]
ExecStart=/usr/bin/mdadm --monitor --scan

[Install]
WantedBy=multi-user.target
