# repos/testing-i686/php-fpm.service
[Unit]
Description=FastCGI Process Manager for PHP
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/run/php-fpm/php-fpm.pid
PrivateTmp=true
ExecStart=/usr/sbin/php-fpm
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
