# repos/testing-i686/slapd.service
[Unit]
Description=OpenLDAP server daemon

[Service]
Type=forking
ExecStart=/usr/sbin/slapd

[Install]
WantedBy=multi-user.target
