# repos/testing-i686/updatedb.cron
#!/bin/sh

NICE='nice -n 19'
IONICE='ionice -c 2 -n 7'

exec ${IONICE} ${NICE} updatedb -f proc
