# repos/testing-i686/var-lib-nfs-rpc_pipefs.mount
[Unit]
Description=RPC Pipe File System 
DefaultDependencies=no

[Mount]
What=sunrpc
Where=/var/lib/nfs/rpc_pipefs
Type=rpc_pipefs
