# qemu/trunk/qemu-ga.service
[Unit]
Description=QEMU Guest Agent
ConditionVirtualization=qemu

[Service]
ExecStart=/usr/bin/qemu-ga

[Install]
WantedBy=multi-user.target
