# conntrack-tools/repos/testing-x86_64/conntrackd.service
[Unit]
Description=Conntrack Daemon

[Service]
ExecStart=/usr/bin/conntrackd

[Install]
WantedBy=multi-user.target
