# hylafax/repos/extra-i686/hylafax.cron.daily
#!/bin/sh

/usr/bin/faxqclean
/usr/bin/faxcron -rcv 30
