# ossp/repos/extra-i686/osspd.service
[Unit]
Description=OSS Userspace Bridge

[Service]
ExecStart=/usr/bin/osspd -f

[Install]
WantedBy=multi-user.target
