# gpm/trunk/gpm.service
[Unit]
Description=GPM daemon

[Service]
Type=forking
ExecStart=/usr/sbin/gpm
