# repos/extra-i686/rygel.service
[Unit]
Description=Rygel DLNA Media Server

[Service]
Type=dbus
BusName=org.gnome.Rygel1
ExecStart=/usr/bin/rygel

[Install]
WantedBy=default.target
