# repos/extra-i686/sensord.service
[Unit]
Description=Sensor information logging daemon

[Service]
PIDFile=/var/run/sensord.pid
ExecStart=/usr/bin/sensord

[Install]
WantedBy=multi-user.target
