# metalog/repos/staging-x86_64/metalog.service
[Unit]
Description=System Logger Daemon

[Service]
ExecStart=/usr/bin/metalog

[Install]
WantedBy=multi-user.target
Alias=syslog.service
