# bftpd/repos/testing-x86_64/bftpd.socket
[Unit]
Description=bftpd incoming socket
Conflicts=bftpd.service

[Socket]
ListenStream=21
Accept=yes

[Install]
WantedBy=sockets.target
