# util-linux/repos/core-x86_64/rfkill-unblock_.service
[Unit]
Description=RFKill-Unblock %I
After=rfkill-block@all.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rfkill unblock %I

[Install]
WantedBy=multi-user.target
