# repos/staging-x86_64/krb5-kpropd@.service
[Unit]
Description=Kerberos 5 propagation server
Conflicts=krb5-kpropd.service

[Service]
ExecStart=/usr/bin/kpropd
StandardInput=socket
StandardError=syslog
