# repos/staging-x86_64/krb5-kdc.service
[Unit]
Description=Kerberos 5 KDC

[Service]
ExecStart=/usr/bin/krb5kdc -n
Restart=always

[Install]
WantedBy=multi-user.target
